# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

echo "[services.d] nginx: listening on port" $HTTP_PORT
/usr/sbin/nginx -g "daemon off;"
