# docker/rootfs/etc/services.d/php/run
#!/usr/bin/with-contenv sh

echo "[services.d] php-fpm: starting"
/usr/local/sbin/php-fpm --nodaemonize -R
